# data/supervisor/mymanager_server.conf
[program:mymanager_server]
directory = /var/www/MyManager/server
environment = BEEGO_RUNMODE="prod"
command = /var/www/MyManager/server/mymanager
user = www-data
autostart = true
autorestart = true
startsecs = 5
redirect_stderr = true
stdout_logfile = /var/www/MyManager/data/log/server-std.log

